# oracle-linux-9-u2/etc/systemd/system/sockets.target.wants/virtnetworkd.socket
[Unit]
Description=Libvirt network local socket
Before=virtnetworkd.service


[Socket]
ListenStream=/run/libvirt/virtnetworkd-sock
Service=virtnetworkd.service
SocketMode=0666
RemoveOnStop=yes

[Install]
WantedBy=sockets.target
