# oracle-linux-9-u2/etc/systemd/system/sockets.target.wants/virtsecretd.socket
[Unit]
Description=Libvirt secret local socket
Before=virtsecretd.service


[Socket]
ListenStream=/run/libvirt/virtsecretd-sock
Service=virtsecretd.service
SocketMode=0666
RemoveOnStop=yes

[Install]
WantedBy=sockets.target
