# oracle-linux-9-u2/etc/systemd/system/sockets.target.wants/virtstoraged.socket
[Unit]
Description=Libvirt storage local socket
Before=virtstoraged.service


[Socket]
ListenStream=/run/libvirt/virtstoraged-sock
Service=virtstoraged.service
SocketMode=0666
RemoveOnStop=yes

[Install]
WantedBy=sockets.target
